# root/etc/s6-overlay/s6-rc.d/init-docker-nginx-crowdsec/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
set -Eeuo pipefail

# Keep generated bouncer credentials private to the container user.
umask 077

API_KEY="${CROWDSEC_NGINX_API_KEY:-}"
BOUNCER_CONF=/run/crowdsec/crowdsec-nginx-bouncer.conf
BOUNCER_CONF_TEMPLATE=/etc/crowdsec/bouncers/crowdsec-nginx-bouncer.conf.template
NGINX_CROWDSEC_INCLUDE=/run/nginx/http.d/crowdsec.conf
NGINX_CROWDSEC_INCLUDE_DEFAULT=/defaults/runtime/nginx/crowdsec.conf

mkdir -p /run/crowdsec /run/nginx/http.d
rm -f "${NGINX_CROWDSEC_INCLUDE}" "${BOUNCER_CONF}"

if [[ -n "${API_KEY}" || -n "${CROWDSEC_LAPI_URL:-}" || -n "${CROWDSEC_APPSEC_URL:-}" ]]; then
    if [[ -z "${API_KEY}" || -z "${CROWDSEC_LAPI_URL:-}" ]]; then
        echo "ERROR: CROWDSEC_NGINX_API_KEY and CROWDSEC_LAPI_URL must both be set when CrowdSec is configured." >&2
        exit 1
    fi
    if [[ -e /config/nginx/http.d/crowdsec.conf ]]; then
        echo "ERROR: Remove the obsolete /config/nginx/http.d/crowdsec.conf before enabling CrowdSec." >&2
        exit 1
    fi

    echo "Configuring CrowdSec nginx bouncer..."

    # The nginx include is static. The bouncer config needs runtime secrets.
    awk \
        -v api_url="${CROWDSEC_LAPI_URL}" \
        -v api_key="${API_KEY}" \
        -v appsec_url="${CROWDSEC_APPSEC_URL:-}" \
        'BEGIN { appsec_written = 0 }
         /^API_URL=/ { print "API_URL=" api_url; next }
         /^API_KEY=/ { print "API_KEY=" api_key; next }
         /^CAPTCHA_PROVIDER=/ { print "CAPTCHA_PROVIDER="; next }
         /^MODE=/ { print "MODE=stream"; next }
         /^REQUEST_TIMEOUT=/ { print "REQUEST_TIMEOUT=500"; next }
         /^APPSEC_URL=/ {
             appsec_written = 1
             if (appsec_url != "") {
                 print "APPSEC_URL=" appsec_url
             } else {
                 print
             }
             next
         }
         { print }
         END {
             if (appsec_url != "" && appsec_written == 0) {
                 print "APPSEC_URL=" appsec_url
             }
         }' \
        "${BOUNCER_CONF_TEMPLATE}" \
        > "${BOUNCER_CONF}"

    if [[ -z ${LSIO_NON_ROOT_USER:-} ]]; then
        chown abc:abc "${BOUNCER_CONF}"
    fi
    chmod 0600 "${BOUNCER_CONF}"
    install -m 0644 "${NGINX_CROWDSEC_INCLUDE_DEFAULT}" "${NGINX_CROWDSEC_INCLUDE}"
else
    echo "CrowdSec configuration not set, nginx bouncer disabled."
fi
