# root/etc/s6-overlay/s6-rc.d/init-docker-nginx-crowdsec/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
set -Eeuo pipefail

# Keep generated bouncer credentials private to the container user.
umask 077

API_KEY="${CROWDSEC_NGINX_API_KEY:-}"
BOUNCER_CONF=/run/crowdsec/crowdsec-nginx-bouncer.conf
NGINX_CROWDSEC_INCLUDE=/run/nginx/http.d/crowdsec.conf
NGINX_CROWDSEC_INCLUDE_DEFAULT=/defaults/runtime/nginx/crowdsec.conf

mkdir -p /run/crowdsec /run/nginx/http.d
rm -f "${NGINX_CROWDSEC_INCLUDE}" "${BOUNCER_CONF}"

if [[ -n "${API_KEY}" || -n "${CROWDSEC_LAPI_URL:-}" || -n "${CROWDSEC_APPSEC_URL:-}" ]]; then
    if [[ -z "${API_KEY}" || -z "${CROWDSEC_LAPI_URL:-}" ]]; then
        echo "ERROR: CROWDSEC_NGINX_API_KEY and CROWDSEC_LAPI_URL must both be set when CrowdSec is configured." >&2
        exit 1
    fi
    if [[ -e /config/nginx/http.d/crowdsec.conf ]]; then
        echo "ERROR: Remove the obsolete /config/nginx/http.d/crowdsec.conf before enabling CrowdSec." >&2
        exit 1
    fi

    echo "Configuring CrowdSec nginx bouncer..."

    # Generate the exact supported configuration rather than rewriting the
    # formatting and defaults of an upstream example file.
    cat > "${BOUNCER_CONF}" <<EOF
ENABLED=true
API_URL=${CROWDSEC_LAPI_URL}
API_KEY=${API_KEY}
USE_TLS_AUTH=false
TLS_CLIENT_CERT=
TLS_CLIENT_KEY=
CACHE_EXPIRATION=1
BOUNCING_ON_TYPE=all
FALLBACK_REMEDIATION=ban
REQUEST_TIMEOUT=500
UPDATE_FREQUENCY=10
ENABLE_INTERNAL=false
MODE=stream
EXCLUDE_LOCATION=
BAN_TEMPLATE_PATH=/var/lib/crowdsec/lua/templates/ban.html
REDIRECT_LOCATION=
RET_CODE=
CAPTCHA_PROVIDER=
SECRET_KEY=
SITE_KEY=
CAPTCHA_TEMPLATE_PATH=
CAPTCHA_EXPIRATION=3600
APPSEC_URL=${CROWDSEC_APPSEC_URL:-}
APPSEC_FAILURE_ACTION=passthrough
APPSEC_CONNECT_TIMEOUT=
APPSEC_SEND_TIMEOUT=
APPSEC_PROCESS_TIMEOUT=
ALWAYS_SEND_TO_APPSEC=false
APPSEC_DROP_UNREADABLE_BODY=false
SSL_VERIFY=true
EOF

    if [[ -z ${LSIO_NON_ROOT_USER:-} ]]; then
        chown abc:abc "${BOUNCER_CONF}"
    fi
    chmod 0600 "${BOUNCER_CONF}"
    install -m 0644 "${NGINX_CROWDSEC_INCLUDE_DEFAULT}" "${NGINX_CROWDSEC_INCLUDE}"
else
    echo "CrowdSec configuration not set, nginx bouncer disabled."
fi
