# root/etc/s6-overlay/s6-rc.d/init-docker-nginx-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
set -Eeuo pipefail

# Copy active default config files only when they do not already exist. The
# adjacent .sample files are refreshed separately by the samples stage, while
# user-managed active files are never replaced.
while IFS= read -r -d '' source; do
    relative_path="${source#/defaults/}"
    destination="/config/${relative_path}"
    if [ ! -e "${destination}" ]; then
        mkdir -p "$(dirname "${destination}")"
        cp "${source}" "${destination}"
    fi
done < <(find /defaults/nginx -type f -print0)
