# root/etc/s6-overlay/s6-rc.d/init-docker-nginx-geoipupdate/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
set -Eeuo pipefail

# Protect the MaxMind license key written to the runtime filesystem.
umask 077

ACCOUNT_ID="${GEOIPUPDATE_ACCOUNT_ID:-}"
LICENSE_KEY="${GEOIPUPDATE_LICENSE_KEY:-}"
EDITION_IDS="${GEOIPUPDATE_EDITION_IDS}"
GEOIP_CONF=/run/GeoIP.conf
BOOTSTRAP_MARKER=/run/geoipupdate-bootstrapped

rm -f "${BOOTSTRAP_MARKER}"

if [[ -n "${ACCOUNT_ID}" || -n "${LICENSE_KEY}" ]]; then
    if [[ -z "${ACCOUNT_ID}" || -z "${LICENSE_KEY}" ]]; then
        echo "ERROR: GEOIPUPDATE_ACCOUNT_ID and GEOIPUPDATE_LICENSE_KEY must both be set when GeoIPUpdate is configured." >&2
        exit 1
    fi
fi

if [[ -n "${ACCOUNT_ID}" && -n "${LICENSE_KEY}" ]]; then
    echo "Configuring GeoIPUpdate..."
    mkdir -p /config/geoip

    cat > "${GEOIP_CONF}" <<EOF
AccountID ${ACCOUNT_ID}
LicenseKey ${LICENSE_KEY}
EditionIDs ${EDITION_IDS}
EOF
    chmod 0600 "${GEOIP_CONF}"

    missing_database=false
    read -r -a edition_ids <<< "${EDITION_IDS}"
    for edition_id in "${edition_ids[@]}"; do
        if [[ ! -f "/config/geoip/${edition_id}.mmdb" ]]; then
            missing_database=true
            break
        fi
    done

    if [[ "${missing_database}" == true ]]; then
        echo "Bootstrapping missing GeoIP database..."
        if ! /usr/local/bin/update-geoip-database; then
            echo "ERROR: Unable to bootstrap the configured GeoIP database." >&2
            exit 1
        fi
        : > "${BOOTSTRAP_MARKER}"
    fi

    echo "GeoIPUpdate configured; the supervised updater will refresh databases."
else
    rm -f "${GEOIP_CONF}"
    echo "GeoIPUpdate credentials not set; periodic updates are disabled."
fi
