# root/etc/s6-overlay/s6-rc.d/init-docker-nginx-geoipupdate/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
set -Eeuo pipefail

# Protect the MaxMind license key written to the runtime filesystem.
umask 077

ACCOUNT_ID="${GEOIPUPDATE_ACCOUNT_ID:-}"
LICENSE_KEY="${GEOIPUPDATE_LICENSE_KEY:-}"
EDITION_IDS="${GEOIPUPDATE_EDITION_IDS}"
GEOIP_CONF=/run/GeoIP.conf

if [[ -n "${ACCOUNT_ID}" || -n "${LICENSE_KEY}" ]]; then
    if [[ -z "${ACCOUNT_ID}" || -z "${LICENSE_KEY}" ]]; then
        echo "ERROR: GEOIPUPDATE_ACCOUNT_ID and GEOIPUPDATE_LICENSE_KEY must both be set when GeoIPUpdate is configured." >&2
        exit 1
    fi
fi

if [[ -n "${ACCOUNT_ID}" && -n "${LICENSE_KEY}" ]]; then
    echo "Configuring GeoIPUpdate..."
    mkdir -p /config/geoip

    cat > "${GEOIP_CONF}" <<EOF
AccountID ${ACCOUNT_ID}
LicenseKey ${LICENSE_KEY}
EditionIDs ${EDITION_IDS}
EOF
    chmod 0600 "${GEOIP_CONF}"

    echo "Running geoipupdate..."
    if ! /usr/local/bin/geoipupdate -f "${GEOIP_CONF}" -d /config/geoip; then
        if compgen -G '/config/geoip/*.mmdb' > /dev/null; then
            echo "WARNING: geoipupdate failed, but an existing database was found. Continuing with cached database."
        else
            echo "ERROR: geoipupdate failed and no existing database was found." >&2
            exit 1
        fi
    fi
else
    rm -f "${GEOIP_CONF}"
    echo "GeoIPUpdate credentials not set, skipping update."
fi
