# root/etc/periodic/daily/geoipupdate
#!/bin/sh
/usr/local/bin/geoipupdate -f /etc/GeoIP.conf -d /config/geoip
