# root/etc/s6-overlay/s6-rc.d/init-docker-nginx-end/up
# This file is the dependency boundary between docker-nginx and the base init process.
