# root/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# Copy active default config files only when they do not already exist. The
# adjacent .sample files are refreshed separately by init-samples, while
# user-managed active files are never replaced.
while IFS= read -r -d '' source; do
    relative_path="${source#/defaults/}"
    # resolver.conf is generated from /etc/resolv.conf by init-resolver. Keep
    # the shipped copy only as the adjacent host-visible sample.
    if [ "${relative_path}" = "nginx/snippets/resolver.conf" ]; then
        continue
    fi
    destination="/config/${relative_path}"
    if [ ! -e "${destination}" ]; then
        mkdir -p "$(dirname "${destination}")"
        cp "${source}" "${destination}"
    fi
done < <(find /defaults/nginx -type f -print0)
