# root/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# Copy active default config files only when they do not already exist. The
# adjacent .sample files are refreshed separately by init-samples, while
# user-managed active files are never replaced.
while IFS= read -r -d '' source; do
    relative_path="${source#/defaults/}"
    destination="/config/${relative_path}"
    if [ ! -e "${destination}" ]; then
        mkdir -p "$(dirname "${destination}")"
        cp "${source}" "${destination}"
    fi
done < <(find /defaults/nginx -type f -print0)

echo "Validating nginx configuration..."
/usr/sbin/nginx -t -e stderr 2>&1
