# root/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# Update default config files
cp -ru /defaults/nginx/ /config/

echo "Validating nginx configuration..."
/usr/sbin/nginx -t
