# root/etc/s6-overlay/s6-rc.d/init-geoipupdate/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
set -e

# Protect the MaxMind license key written to the runtime filesystem.
umask 077

startup_error() {
    echo "ERROR: $*" >&2
    printf '1\n' > /run/s6-linux-init-container-results/exitcode
    /run/s6/basedir/bin/halt
    return 1
}

ACCOUNT_ID="${GEOIPUPDATE_ACCOUNT_ID:-}"
LICENSE_KEY="${GEOIPUPDATE_LICENSE_KEY:-}"
EDITION_IDS="${GEOIPUPDATE_EDITION_IDS}"
GEOIP_CONF=/run/GeoIP.conf

if [[ -n "${ACCOUNT_ID}" || -n "${LICENSE_KEY}" ]]; then
    if [[ -z "${ACCOUNT_ID}" || -z "${LICENSE_KEY}" ]]; then
        startup_error \
            "GEOIPUPDATE_ACCOUNT_ID and GEOIPUPDATE_LICENSE_KEY must both be set when GeoIPUpdate is configured."
        exit 1
    fi
fi

if [[ -n "${ACCOUNT_ID}" && -n "${LICENSE_KEY}" ]]; then
    echo "Configuring GeoIPUpdate..."
    mkdir -p /config/geoip

    cat > "${GEOIP_CONF}" <<EOF
AccountID ${ACCOUNT_ID}
LicenseKey ${LICENSE_KEY}
EditionIDs ${EDITION_IDS}
EOF
    chmod 0600 "${GEOIP_CONF}"

    echo "Running geoipupdate..."
    if ! /usr/local/bin/geoipupdate -f "${GEOIP_CONF}" -d /config/geoip; then
        if compgen -G '/config/geoip/*.mmdb' > /dev/null; then
            echo "WARNING: geoipupdate failed, but an existing database was found. Continuing with cached database."
        else
            startup_error "geoipupdate failed and no existing database was found."
            exit 1
        fi
    fi
else
    rm -f "${GEOIP_CONF}"
    echo "GeoIPUpdate credentials not set, skipping update."
fi
